# weather/weather-anomaly.service
[Unit]
Description=Solar PV anomaly detector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/chris/weather
EnvironmentFile=/home/chris/weather/weather.env
ExecStart=/usr/bin/python3 /home/chris/weather/anomaly.py
User=chris

[Install]
WantedBy=multi-user.target
